# service/rlsudo.service
[Unit]
Description=rlsudo

[Service]
ExecStart=/bin/bash /opt/rlsudo/start_rlsudo.sh

[Install]
WantedBy=multi-user.target