# service/rludo.service
[Unit]
Description=rlsudo

[Service]
ExecStart=/bin/bash /opt/start_rlsudo.sh

[Install]
WantedBy=multi-user.target